# mbpro/services/systemd-suspend.service.d/lid.conf
[Service]
ExecStartPre=/usr/local/bin/lid.sh prepare
ExecStartPost=/usr/local/bin/lid.sh restore
